# nellx.service
[Unit]
Description=NellX Telegram Bot and Web App
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/obmen
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/python3 main.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
